# ci/1030-docker-run/etc/profile.d/sm3sum.sh
alias sm3sum 'cksum -a sm3' 2>/dev/null
